# molecule/timer_service_removal/files/test.service
[Unit]
Description=test timer service
Requires=test.timer

[Service]
Type=simple
ExecStart=echo "this is a test service"
User=root
Group=root

