# molecule/timer_service/files/test.timer
[Unit]
Description=test timer

[Timer]
OnBootSec=1year
Unit=test.service

[Install]
WantedBy=timers.target
